# versal/sw/yocto-layers/meta-vssr-trd/recipes-utils/watchdog-daemon-init/files/watchdog-daemon-init.service
[Unit]
Description=watchdog-daemon-init

[Service]
ExecStart=/usr/bin/watchdog-daemon

[Install]
WantedBy=multi-user.target

